# crontab
# m h dom mon dow command
30 * * * * update-page >> /var/log/update-pages.log 2>&1

